# packaging/systemd/dimeng-monitor-agent.service
[Unit]
Description=DiMeng Monitor Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dimeng-agent
Group=dimeng-agent
EnvironmentFile=/etc/dimeng-monitor-agent/agent.env
ExecStart=/usr/local/bin/dimeng-monitor-agent --endpoint ${DIMENG_ENDPOINT} --claim-token-file /var/lib/dimeng-monitor-agent/claim.token
Restart=always
RestartSec=10
UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectHome=yes
ProtectSystem=strict
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
CapabilityBoundingSet=
AmbientCapabilities=
ReadWritePaths=/var/lib/dimeng-monitor-agent
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
MemoryMax=64M
CPUQuota=20%

[Install]
WantedBy=multi-user.target
